# supervisor.conf
[program:rfs]
directory=/home/rfs
command=node /home/rfs
stdout_logfile=/var/rfs/logs/supervisor.log
redirect_stderr=true
